# roles/monitoring/templates/prometheus.service.j2
[Unit]
Description=Prometheus
After=network.target

[Service]
ExecStart=/usr/local/bin/prometheus

[Install]
WantedBy=multi-user.target